# misc/supervisor.conf
[program:uwsgi]
command=/home/showthedocs/env/bin/uwsgi
  --ini /home/showthedocs/repo/misc/uwsgi.ini
autostart=true
autorestart=true
user=showthedocs
